# systemd/ratools-rad.socket
[Unit]
Description=Socket for IPv6 Router Advertisement Daemon

[Socket]
ListenStream=/var/run/ratools/rad.sock

[Install]
WantedBy=sockets.target
